# ansiblecube/roles/system/files/unblock-rf.service
[Unit]
Description=Ublock RF devices
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin rfkill ublock all
