# ansiblecube/roles/resize/templates/mount-boot-partition.service
[Unit]
Description=Mount Boot Partition (custom)
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/mount /dev/mmcblk0p1
